# DGTCentaurMods/DEBIAN/postinst
#!/bin/bash

ARCH=`uname -m`
if [ -f /boot/firmware/config.txt ]
then
    CONFIG="/boot/firmware/config.txt"
else
    CONFIG="/boot/config.txt"
fi
PACKAGE="DGTCentaurMods"
DGTCM_PATH="/opt/${PACKAGE}"
if [ -f /boot/firmware/cmdline.txt ]
then
    CMDLINEFILE="/boot/firmware/cmdline.txt"
else
    CMDLINEFILE="/boot/cmdline.txt"
fi
CENTAUR="/home/pi/centaur/"
CENTAUR_ENGINE="/home/pi/centaur/engines"
STOCKFISH="stockfish_pi"
RODENT_REPO="https://github.com/nescitus/rodent-iv.git"
MAIA_REPO="https://github.com/CSSLab/maia-chess.git"
FENLOG="/home/pi/centaur/fen.log"
VENV_PATH="${DGTCM_PATH}/.venv"


if ! cat $CONFIG | grep -q "#dgtcentaurmods"
then
    echo -e "::: Hardware configuration."
    echo -e "::: Backup few files."
    cp $CONFIG $CONFIG.bak
    cp $CMDLINEFILE $CMDLINEFILE.bak

    # Enable SPI bus if not enabled
    echo "::: Checking SPI"
        SPION=`cat $CONFIG | grep "^#dtparam=spi=on"`
        SPI10=`cat $CONFIG | grep "^dtoverlay=spi1-3cs"`
    if [ ! -z $SPION ]
    then
        sed -i "s/#dtparam=spi=on/dtparam=spi=on/g" $CONFIG
    else
        sed -i "$ a dtparam=spi=on" $CONFIG
    fi    
    
    echo -e "::: Checking SPI 1.0 bus."
    if [ -z $SPI10 ]
    then
        sed -i "$ a dtoverlay=spi1-3cs" $CONFIG
    fi

echo "::: Checking serial port."
    UARTON=`cat $CONFIG | grep "^#enable_uart=1"`
    if [ ! -z $UARTON ]
    then
        echo "::: Enabling serial port."
        sed -i "s/#enable_uart=1.*/enable_uart=1/g" $CONFIG
    else
        sed -i "$ a enable_uart=1" $CONFIG
    fi
    # Mark file as modified
    sed -i "$ a #dgtcentaurmods" $CONFIG
fi

if cat $CMDLINEFILE | grep -q "console=serial0"
then
    echo -e "::: Disable console on ttyS0"
        REPLY=$(sed 's/[^ ]* *//' $CMDLINEFILE)
        echo -e "$REPLY" > $CMDLINEFILE
fi

# Make DGTCM available to all users
rm /usr/lib/python3/dist-packages/${PACKAGE}
cd /usr/lib/python3/dist-packages && sudo ln -s ${DGTCM_PATH} ${PACKAGE}

# Allow flask on port 80
sudo touch /etc/authbind/byport/80
sudo chmod 777 /etc/authbind/byport/80

function enableTimeSyncService {
    # Keep system time sane for TLS/API interactions without forcing a specific daemon.
    for service in systemd-timesyncd ntpsec ntp chrony
    do
        if systemctl list-unit-files "${service}.service" 2>/dev/null | grep -q "^${service}.service"
        then
            echo -e "::: Enabling time sync service: ${service}."
            systemctl enable "${service}.service" >/dev/null 2>&1 || true
            return
        fi
    done

    echo -e "::: No known time sync service found. Skipping enable step."
}

# Enable systemd services.
echo -e "::: Setup services."
    systemctl daemon-reload
    systemctl enable ${PACKAGE}Web.service
    systemctl enable ${PACKAGE}Update.service
    systemctl enable ${PACKAGE}.service
    enableTimeSyncService

systemctl daemon-reload

# Check if CENTAUR directory exists
if [ ! -d $CENTAUR ]
then
    mkdir $CENTAUR
fi

# Check if CENTAUR engine directory exists
if [ ! -d $CENTAUR_ENGINE ]
then
    mkdir $CENTAUR_ENGINE
fi

# Check if fen.log exists
if [ ! -e $FENLOG ]
then
    touch $FENLOG
fi

echo -e "::: Creating python virtual environment."
    cd $DGTCM_PATH

    # Always recreate venv to avoid stale pip packages surviving upgrades.
    if [ -d "${VENV_PATH}" ]
    then
        rm -rf "${VENV_PATH}"
    fi

    python3 -m venv --system-site-packages "${VENV_PATH}"
    if [ $? -ne 0 ]
    then
        echo "::: ERROR: Failed to create virtual environment at ${VENV_PATH}."
        exit 1
    fi

echo -e "::: Installing pip-only python packages into virtual environment."
    "${VENV_PATH}/bin/python" -m pip install --upgrade pip setuptools wheel
    CFLAGS="-fcommon" "${VENV_PATH}/bin/pip" install -r $DGTCM_PATH/config/pip-requirements.txt
    if [ $? -ne 0 ]
    then
        echo "::: ERROR: Python dependency installation failed in virtual environment."
        exit 1
    fi


function insertMaia {
    if [[ $ARCH == "armv6"* ]]; then
        #Remove maia and don't get weights
        echo -e "::: This system is not compatible with maia engine."
        rm -rf ${DGTCM_PATH}/engines/maia* > /dev/null 2>&1
        return
    fi
    if [ ! -d  ${DGTCM_PATH}/tmp ]
    then
        mkdir ${DGTCM_PATH}/tmp
    fi
    if [ ! -d ${DGTCM_PATH}/tmp/maia-chess ]
    then
        git clone --depth 1 $MAIA_REPO ${DGTCM_PATH}/tmp/maia-chess;
        if (($? > 0))
            then echo "git clone maia chess failed with $?"
        fi
    fi
    cp -rn ${DGTCM_PATH}/tmp/maia-chess/maia_weights ${DGTCM_PATH}/engines/
    chmod 755 ${DGTCM_PATH}/engines/maia
}

function insertRodentIV {
    if [ ! -d  ${DGTCM_PATH}/tmp ]
    then
        mkdir ${DGTCM_PATH}/tmp
    fi
    if [ ! -d ${DGTCM_PATH}/tmp/rodent-iv ]
    then
        git clone --depth 1 $RODENT_REPO ${DGTCM_PATH}/tmp/rodent-iv;
        if (($? > 0))
            then echo "git clone rodentiv failed with $?"
        fi
    fi
    cp -rn ${DGTCM_PATH}/tmp/rodent-iv/personalities ${DGTCM_PATH}/engines/
    cp -rn ${DGTCM_PATH}/tmp/rodent-iv/books ${DGTCM_PATH}/engines/
    chmod 744 ${DGTCM_PATH}/engines/rodentIV
}

function insertZahak {
    chmod 755 ${DGTCM_PATH}/engines/zahak
}

insertMaia
insertRodentIV
insertZahak

# Setting final permissions
function setPermissions {
    echo -e "::: Setting permissions."
    sudo chown -R pi.pi ${DGTCM_PATH}
    sudo chown -R pi.root $CENTAUR

}

function copyDefaults {
    echo -e "::: Copying default values."
    cp -rn ${DGTCM_PATH}/defaults/* ${DGTCM_PATH}/
}

case $1 in
    configure* )

    copyDefaults
    setPermissions

    echo -e "::: Done."
    ;;
esac

exit 0

